# deploy/baseline-break-monitor.service
[Unit]
Description=Baseline Break tennis live monitor (re-runs the model)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/tennis-forecast
Environment=PYTHONUNBUFFERED=1
# Kalshi creds piggyback on the gas-bot's .env (same key reused
# across all bots in the suite, read-only). The leading ``-`` makes
# systemd skip the line if the file doesn't exist on a fresh droplet.
EnvironmentFile=-/root/gas-prices/.env
EnvironmentFile=-/root/tennis-forecast/.env
ExecStart=/root/tennis-forecast/.venv/bin/python /root/tennis-forecast/scripts/run_live_monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
